# profile.d/path.sh
# Add ~/bin folder to path
export PATH=~/bin:$PATH

# Setting PATH for Python 3.4
export PATH="/Library/Frameworks/Python.framework/Versions/3.4/bin:${PATH}"
